# window-renicer.service
[Service]
ExecStart=/usr/bin/window-renicer

[Install]
WantedBy=multi-user.target
