# window-renicer.service
[Service]
ExecStart=/usr/bin/window-renicer

[Install]
WantedBy=default.target
